# systemd/system/refresh-policy-routes@.timer
[Timer]

OnActiveSec=30
OnUnitInactiveSec=60
RandomizedDelaySec=5

[Install]
WantedBy=setup-policy-routes@%i.service
